# deploy/zhtp-dev.service
[Unit]
Description=ZHTP Node (Development)
After=network.target

[Service]
Type=simple
ExecStart=/opt/zhtp/zhtp node start --network dev --dev --config /opt/zhtp/config.toml
WorkingDirectory=/opt/zhtp
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=info
Environment=ZHTP_ALLOW_BOOTSTRAP=1

[Install]
WantedBy=multi-user.target
